# packaging/systemd/photoframe.service
[Unit]
Description=Digital Photo Frame Backend
After=network.target
Wants=network.target

[Service]
Type=simple
User=photoframe
Group=photoframe
WorkingDirectory=/etc/photoframe
ExecStart=/usr/bin/digital-photo-frame
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security (relaxed for compatibility)
NoNewPrivileges=yes
PrivateTmp=no
ProtectSystem=no
ProtectHome=no

[Install]
WantedBy=multi-user.target